# web/service.conf
# sigvisa_web.conf
start on filesystem
setuid sigvisa
exec /home/sigvisa/python/sigvisa/runweb.sh